# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:ttyd]
command=/usr/bin/ttyd -p 7681 -W -t fontSize=16 -t enableMouse=true -t allow_osc52=true -t cursorBlink=true -t disableLeaveAlert=true tmux new -A -s main
user=nimda
directory=/home/nimda/workspace
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:code-server]
command=/usr/bin/code-server --bind-addr 0.0.0.0:8080 --auth none /home/nimda/workspace
user=nimda
directory=/home/nimda
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dashboard]
command=/home/nimda/_SPACE/dash/dashboard
user=nimda
directory=/home/nimda/_SPACE/dash
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:spacebot]
command=/home/nimda/scripts/run_spacebot.sh
user=nimda
directory=/home/nimda
environment=HOME="/home/nimda",USER="nimda"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pinchtab]
command=/usr/bin/python3 pinchtab_mcp.py
user=nimda
directory=/home/nimda
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cloudflared]
command=/home/nimda/scripts/run_cloudflared.sh
user=nimda
directory=/home/nimda
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
