# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:ttyd]
command=/usr/bin/ttyd -p 7681 -b /ttyd -W -t fontSize=16 -t enableMouse=true tmux new -A -s main
user=nimda
directory=/home/nimda
autorestart=true

[program:code-server]
command=/usr/bin/code-server --bind-addr 0.0.0.0:8080 --auth none
user=nimda
directory=/home/nimda
autorestart=true

[program:dashboard]
command=/home/nimda/_SPACE/dash/dashboard
user=nimda
directory=/home/nimda/_SPACE/dash
environment=XDG_RUNTIME_DIR="/run/user/1000",XAUTHORITY="/home/nimda/.Xauthority",DISPLAY=":0"
autorestart=true

[program:spacebot]
command=/home/nimda/spacebot start
user=nimda
directory=/home/nimda
environment=HOME="/home/nimda",USER="nimda"
autorestart=true

[program:pinchtab]
command=/usr/bin/python3 pinchtab_mcp.py
user=nimda
directory=/home/nimda
autorestart=true

[program:cloudflared]
command=/usr/local/bin/run_cloudflared.sh
user=nimda
directory=/home/nimda
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
